# .config/systemd/user/pulseaudio.service
# ~/.config/systemd/user/pulseaudio.service

[Unit]
Description=PulseAudio Server
Requires=dbus.socket

[Service]
ExecStart=/usr/bin/pulseaudio
ExecStop=/usr/bin/pulseaudio --kill

BusName=org.pulseaudio.Server
Restart=on-failure

[Install]
WantedBy=default.target
