# .local/share/runit/service/transmission/run
#!/bin/sh --
exec transmission-daemon -f --log-level=error 2>&1
