# local/service/transmission/run
#!/bin/sh
exec transmission-daemon -f --log-error 2>&1
