# .local/share/runit/service/pulseaudio/run
#!/bin/sh --
exec pulseaudio --exit-idle-time=-1 2>&1
