# .local/share/runit/service/mpd/run
#!/bin/sh --
if sv check pulseaudio; then
    exec mpd --no-daemon --stderr 2>&1
fi
