# local/service/mpd/run
#!/bin/sh
set -e
sv check pulseaudio
exec 2>&1 mpd --no-daemon
